# scripts/kanban-sync.service
[Unit]
Description=Kanban sync from GitHub
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/kanban-sync.sh
User=%i
Environment=HOME=/home/%i

[Install]
WantedBy=multi-user.target
